# rootfs/config/supervisor/1-xvfb.conf
[program:Xvfb]
command=bash -c 'Xvfb $DISPLAY -screen 0 $VNC_RESOLUTION"x24" -nolisten unix'
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
priority=1

