# rootfs/config/supervisor/2-openbox.conf
[program:Openbox]
command=openbox

priority=2
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

