# rootfs/config/supervisor/6-self-ping.conf
[program:Self-ping]
command=bash /config/self-ping.sh
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
exitcodes=0
startsecs=0
priority=6
